# ha_opencode/rootfs/etc/s6-overlay/s6-rc.d/codenomad/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# CodeNomad - Cleanup on exit
# ==============================================================================

declare APP_EXIT_CODE=${1}

if [[ "${APP_EXIT_CODE}" -eq 0 ]]; then
    bashio::log.info "CodeNomad stopped cleanly"
elif [[ "${APP_EXIT_CODE}" -eq 256 ]]; then
    # Exit code 256 typically means SIGTERM (normal shutdown)
    bashio::log.info "CodeNomad terminated by signal"
else
    bashio::log.warning "CodeNomad exited with code ${APP_EXIT_CODE}"
fi
