# ha_opencode/rootfs/etc/s6-overlay/s6-rc.d/codenomad/run
#!/command/with-contenv bashio
# ==============================================================================
# CodeNomad - Web frontend for OpenCode server mode
# ==============================================================================

# Check if CodeNomad is enabled
CODENOMAD_ENABLED=$(bashio::config 'codenomad_enabled' 2>/dev/null || echo "false")
if [ "${CODENOMAD_ENABLED}" != "true" ]; then
    bashio::log.info "CodeNomad is disabled, exiting..."
    sleep 10
    exit 0
fi

# Read configuration
CODENOMAD_PORT=$(bashio::config 'codenomad_port' 2>/dev/null || echo "3000")
CODENOMAD_WORKSPACE=$(bashio::config 'codenomad_workspace' 2>/dev/null || echo "/homeassistant")
CODENOMAD_USERNAME=$(bashio::config 'codenomad_username' 2>/dev/null || echo "")
CODENOMAD_PASSWORD=$(bashio::config 'codenomad_password' 2>/dev/null || echo "")

bashio::log.info "Starting CodeNomad on port ${CODENOMAD_PORT}..."
bashio::log.info "Workspace: ${CODENOMAD_WORKSPACE}"

# Set environment
export HOME="/data"
export XDG_DATA_HOME="/data/.local/share"
export XDG_CONFIG_HOME="/data/.config"

# Load user-defined environment variables
if [ -f /data/.env_vars ]; then
    source /data/.env_vars
fi

# Change to workspace directory
cd "${CODENOMAD_WORKSPACE}"

# Set credentials 
export CODENOMAD_SERVER_USERNAME="${CODENOMAD_USERNAME}"
export CODENOMAD_SERVER_PASSWORD="${CODENOMAD_PASSWORD}"

# Start CodeNomad
exec codenomad --http-port ${CODENOMAD_PORT}
